# data/systemd-overrides/xdg-permission-store.service.d/50-sailfishos.conf
# The upstream service uses PartOf=graphical-service.target
# which is not in use on Sailfish OS.
# So we override the dependency here.
[Unit]
PartOf=user-session.target

