# data/xdg-desktop-portal-sailfish.service.in
[Unit]
Description=Portal service for Sailfish OS
PartOf=user-session.target
After=lipstick.service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/xdg-desktop-portal-sailfish
BusName=org.freedesktop.impl.portal.desktop.sailfish
Slice=session.slice
# Auto restart is off as we are DBus activated
Restart=no
